# supervisord-transmission.ini
[program:torrent]
command=/usr/bin/transmission-daemon -f --config-dir /etc/transmission-daemon/
user=media
autostart=true
autorestart=true
startretries=3
stopsignal=QUIT
stopwaitsecs=10
